# swf-monitor-mcp-asgi.service
[Unit]
Description=SWF Monitor MCP endpoint on ASGI (uvicorn) worker
After=network.target postgresql.service

[Service]
Type=simple
User=wenauseic
Group=eic
WorkingDirectory=/opt/swf-monitor/current/src
EnvironmentFile=/opt/swf-monitor/config/env/production.env
Environment=DJANGO_SETTINGS_MODULE=swf_monitor_project.settings
# uvicorn with 2 workers — one event loop per worker, each multiplexing many
# streaming MCP sessions. Bound to loopback only; Apache ProxyPass fronts it.
ExecStart=/opt/swf-monitor/current/.venv/bin/uvicorn \
    swf_monitor_project.asgi:application \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 2 \
    --log-level info \
    --proxy-headers \
    --forwarded-allow-ips 127.0.0.1
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
